# supervisor/NABS.conf
[program:NABS]
command=/opt/NABS/venv/bin/gunicorn -b localhost:8000 -w 4 NABS:app
directory=/opt/NABS
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true