# supervisor/nabs.conf
[program:app]
command=/opt/NABS/venv/bin/gunicorn -b localhost:8000 -c config_gunicorn.py -w 4 app:app
directory=/opt/NABS
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
